# crawler/crontab
*/10 * * * * /usr/bin/bash /home/app/crawler/run_crawler
